# config/docker/web/base/crontab
MAILTO=""
SHELL=/bin/bash

22 * * * * www-data . /cron-env; if [[ -n "$CRON_OK" ]]; then /bin/bash -c 'drush core-cron'; fi
